# backend/ops/systemd/neurocrop-monitor.service
[Unit]
Description=NeuroCrop platform health monitor
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/neurocrop-monitor.env
ExecStart=/opt/neurocrop-backend/scripts/monitor-platform.sh
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/neurocrop-monitor
