# backend/ops/systemd/neurocrop-backup.service
[Unit]
Description=NeuroCrop PostgreSQL backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/neurocrop-backup.env
ExecStart=/opt/neurocrop-backend/scripts/backup-database.sh
User=root
Group=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/neurocrop -/mnt/neurocrop-offsite /var/lock
